# FW/package/ExtBoard/files/etc/hotplug.d/usb/10-ExtBoard
#!/bin/sh

if [ "${DRIVER}" == "rndis_host" ]
then
	#扫描RNDIS设备
	logger -t hotplug RNDIS Device Scan !
	[ -x /usr/bin/ExtBoardRNDISScan ] && /usr/bin/ExtBoardRNDISScan
fi
